# deployment/systemd/controller.service
[Unit]
Description=FEMS Controller Service (development)
After=network.target

[Service]
Type=simple
# Adjust WorkingDirectory to the repository path where `simulator/controller_app.py` lives.
WorkingDirectory=/path/to/fems
ExecStart=/usr/bin/python3 simulator/controller_app.py
Restart=on-failure
User=%i
Environment=PYTHONUNBUFFERED=1
# Ensure controller uses the SQLite event DB; replace the path with the
# actual installation path before enabling the service.
Environment=EVENT_DB=/path/to/fems/simulator/events.db

[Install]
WantedBy=multi-user.target

# Note: This is a template file. Replace `WorkingDirectory` and `User` before enabling.
